# backend/systemd/grpcconnector.service
[Unit]
Description=Zero-Trust gRPC Connector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/grpcconnector run
Restart=always
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
MemoryDenyWriteExecute=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallFilter=@system-service
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
