# systemd/grpctunneler.service
[Unit]
Description=Zero-Trust gRPC Tunneler
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=/etc/grpctunneler/tunneler.conf
# Controller CA is provided via systemd credentials (test/dev only)
LoadCredential=CONTROLLER_CA:/etc/grpctunneler/ca.crt

# Restart behavior (same semantics as connector)
Restart=always
RestartSec=10

# Start command
ExecStart=/usr/bin/grpctunneler run

# Identity & isolation
DynamicUser=yes
PrivateUsers=yes
PrivateTmp=yes
PrivateDevices=yes
NoNewPrivileges=yes

# Filesystem & kernel hardening
ProtectSystem=full
ProtectHome=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes

# Capability & device lockdown
CapabilityBoundingSet=
DevicePolicy=closed

# Syscall & namespace restrictions
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
SystemCallArchitectures=native
SystemCallFilter=~@clock @debug @module @mount @raw-io @reboot @privileged @resources @cpu-emulation @obsolete

# Memory protections
MemoryDenyWriteExecute=yes
LockPersonality=yes

# Runtime behavior
UMask=0777
RuntimeDirectory=grpctunneler

[Install]
WantedBy=multi-user.target
